# infra/systemd/guestbook-api.service
[Unit]
Description=Guestbook submission API for florianvasin.com
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node /opt/guestbook-api/server.js
WorkingDirectory=/opt/guestbook-api
Environment=GB_PENDING=/var/www/florianvasin.com/guestbook-pending.json
Environment=GB_PORT=8041
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
